# deploy/podman/fizzy.service
[Unit]
Description=Fizzy
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/fizzy
ExecStartPre=-/usr/bin/podman rm -f fizzy
ExecStart=/usr/bin/podman run \
  --name fizzy \
  --pull=never \
  --env-file /opt/fizzy/.env \
  --publish 127.0.0.1:3006:80 \
  --volume /srv/fizzy/storage:/rails/storage:Z \
  --health-cmd "curl --fail --silent --show-error --output /dev/null http://127.0.0.1/up" \
  --health-interval=30s \
  --health-timeout=5s \
  --health-retries=3 \
  --health-start-period=60s \
  --health-on-failure=restart \
  localhost/fizzy:local
ExecStop=/usr/bin/podman stop -t 30 fizzy
ExecStopPost=-/usr/bin/podman rm -f fizzy
Restart=always
RestartSec=5
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
